# packaging/systemd/clipvaultd.service
[Unit]
Description=Clipvault clipboard daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/clipvaultd watch
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
